# ansible/guhcampos/podman/roles/podman-redis/templates/podman-redis.service.j2
[Unit]
Description=podman redis
After=podman.service
Requires=podman.service

[Install]
WantedBy=multi-user.target

[Service]
RemainAfterExit=true
WorkingDirectory={{redis_home}}
ExecStart=/usr/bin/podman-compose up -d --remove-orphans
ExecStop=/usr/bin/podman-compose down
