# ansible/guhcampos/podman/roles/mariadb/templates/podman-mariadb.service.j2
[Unit]
Description=podman mariadb
After=podman.service
Requires=podman.service

[Install]
WantedBy=multi-user.target

[Service]
RemainAfterExit=true
WorkingDirectory={{mariadb_home}}
ExecStart=/usr/bin/podman compose up --remove-orphans
ExecStop=/usr/bin/podman compose down
